# crontab
*/1 * * * * docker restart fastapi > /dev/stdout
