# server/cronjob/crontab
*/1 * * * * python /cronjob/cronjob.py >> /cronjob/logfile
  